# contrib/packages/archlinux/nodejs-galilel-explorer-sync-coin.service
[Unit]
Description = Diaz Explorer (Sync) - Coin information update

[Service]
Type = simple
WorkingDirectory = /usr/lib/node_modules/diaz-explorer
ExecStart = /usr/bin/node cron/coin.js
SyslogIdentifier = diaz-explorer
StandardOutput = syslog
StandardError = syslog
