# contrib/packages/archlinux/nodejs-galilel-explorer-sync-proposal.timer
[Unit]
Description = Diaz Explorer (Sync) - Proposal update timer

[Timer]
OnCalendar = *:0/5
RandomizedDelaySec = 5
Persistent = true
 
[Install]
WantedBy = timers.target
