# contrib/packages/archlinux/nodejs-galilel-explorer-sync-proposal.service
[Unit]
Description = Diaz Explorer (Sync) - Proposal update

[Service]
Type = simple
WorkingDirectory = /usr/lib/node_modules/diaz-explorer
ExecStart = /usr/bin/node cron/proposal.js
SyslogIdentifier = diaz-explorer
StandardOutput = syslog
StandardError = syslog
