# contrib/packages/archlinux/nodejs-galilel-explorer-sync-coin.timer
[Unit]
Description = Diaz Explorer (Sync) - Coin information update timer

[Timer]
OnCalendar = *:0/5
RandomizedDelaySec = 30
Persistent = true
 
[Install]
WantedBy = timers.target
