# contrib/packages/archlinux/nodejs-galilel-explorer-sync-peer.timer
[Unit]
Description = Diaz Explorer (Sync) - Network peers update timer

[Timer]
OnCalendar = *:0/1
RandomizedDelaySec = 5
Persistent = true
 
[Install]
WantedBy = timers.target
