# contrib/packages/archlinux/nodejs-galilel-explorer-sync-masternode.service
[Unit]
Description = Diaz Explorer (Sync) - Masternode update

[Service]
Type = simple
WorkingDirectory = /usr/lib/node_modules/diaz-explorer
ExecStart = /usr/bin/node cron/masternode.js
SyslogIdentifier = diaz-explorer
StandardOutput = syslog
StandardError = syslog
